# host/certificates/templates/wapps-acme-renew.timer.j2
[Unit]
Description=Run {{ platform_slug }} ACME renewal daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target

